# templates/kibana.service.erb
# THIS FILE IS MANAGED BY PUPPET!
[Unit]
Description=Kibana daemon
After=local-fs.target network.target

[Service]
User=etcd
Type=simple
WorkingDirectory=/services/kibana/current
ExecStart=/services/kibana/current/bin/kibana 
PIDFile=/services/kibana/kibana.pid

[Install]
WantedBy=multi-user.target
